# root/etc/services.d/duplicacy/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

umask "${UMASK}"

# shellcheck disable=SC2086
HOME="${APP_DIR}" exec s6-setuidgid hotio "${APP_DIR}/duplicacy_web" ${ARGS}
